# upstart/joola.conf
description "joola"

start on started networking
stop on shutdown

chdir /opt/joola
setuid vagrant
setgid vagrant

script
  exec npm start
end script