# upstart/joola.conf
description "joola"

start on started networking
stop on shutdown

chdir /opt/joola
setuid vagrant
setgid vagrant

respawn

script
  exec npm start
end script